# fbcp.service
[Unit]
Description=fbcp

[Service]
ExecStart=/usr/local/bin/fbcp
Type=simple

[Install]
WantedBy=default.target